# RapidBuild64/001-core/RLBFILES/rootcopy/etc/profile.d/rl_numcpus.sh
#!/bin/bash

NUMCPUS=`lscpu | grep 'CPU(s):' | grep -v NUMA | awk '{print $2}'`
export MAKEALLCPUS="-j$NUMCPUS"
